# supervisor/nmbd.conf
[program:nmbd]
command=/usr/sbin/nmbd -F -S
priority=100
autostart=true
autorestart=true
startsecs=5

