# supervisor/qemuviewer.conf
[program:qemuviewer]
command=/usr/local/bin/waitfordaemon qemuconsole /usr/local/bin/qemuviewer
priority=999
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
