# supervisor/qemuconsole.conf
[program:qemuconsole]
command=/usr/local/bin/qemuconsole
priority=500
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
