# supervisor/smbd.conf
[program:smbd]
command=/usr/local/bin/waitfordaemon nmbd /usr/sbin/smbd -F -S
priority=101
autostart=true
autorestart=true
startsecs=5

