# supervisor/smbd.conf
[program:nmbd]
command=/usr/sbin/smbd -F -S
priority=101
autostart=true
autorestart=true
startsecs=5

