# supervisor/xterm.conf
[program:xterm]
command=/usr/bin/xterm
priority=300
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
