# supervisor/xterm.conf
[program:xterm]
command=/usr/local/bin/waitfordaemon vncserver /usr/bin/xterm
priority=300
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=1
startretries=1000000
